# conf/pt_init.conf
start on [2345]
stop on [06]

script
    cd /home/ubuntu/pt
    exec uwsgi --ini /etc/pt/uwsgi.ini
end script
